# qusal/qubes-builder/files/server/rc.local
# SPDX-FileCopyrightText: 2023 The Qubes OS Project <https://www.qubes-os.org>
#
# SPDX-License-Identifier: GPL-3.0-or-later

mount /builder -o dev,suid,remount
