# qusal/sys-ssh/files/server/rc.local
systemctl unmask ssh
systemctl --no-block restart ssh
