# qusal/ansible/files/client/rc.local
systemctl unmask ssh
systemctl --no-block start ssh
