# salt/sys-rsync/files/server/rc.local
systemctl unmask rsync
systemctl --no-block restart rsync
